# s6/root/etc/services.d/ipfs/run
#!/usr/bin/with-contenv bash

# seems like IPFS logs most logs to stderr, so
# map stderr to stdout so s6-log can catch it
exec \
	fdmove -c 2 1 \
	s6-setuidgid nobody ipfs daemon
