# s6/root/etc/services.d/ipfs/run
#!/usr/bin/with-contenv bash

exec \
	fdmove -c 2 1 \
	s6-setuidgid abc ipfs daemon
