# s6/root/etc/services.d/myprog/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid nobody s6-pause
